# launchd/com.moonstone.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.moonstone.daemon</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/moonstone-daemon</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>ThrottleInterval</key>
    <integer>1</integer>

    <key>StandardOutPath</key>
    <string>/var/log/moonstone/daemon.log</string>

    <key>StandardErrorPath</key>
    <string>/var/log/moonstone/daemon.err</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>RUST_LOG</key>
        <string>moonstone=info</string>
    </dict>
</dict>
</plist>
